# services/package-events.service
[Unit]
Description=Update desktop icons

[Service]
ExecStart=/home/user_replace/.config/rofi/launcher/filter.sh
Type=oneshot
User=root
