# services/package-events.path
[Unit]
Description=Monitor package changes

[Path]
PathChanged=/var/lib/dpkg/info/
PathChanged=/var/lib/dpkg/status
PathChanged=/var/log/apt/history.log
Unit=package-events.service

[Install]
WantedBy=multi-user.target
